# output/operational-proof/linux-proof-tmp/clinic-backend-restore.service
[Unit]
Description=Clinic Backend (restore)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/clinic/backend
EnvironmentFile=/root/clinic/backend/.env.restore
Environment=PYTHONPATH=/root/clinic/backend
ExecStart=/root/clinic/backend/.venv/bin/python start_server.py
Restart=always
RestartSec=5
TimeoutStartSec=120
StandardOutput=append:/root/clinic/output/vps/backend-restore.log
StandardError=append:/root/clinic/output/vps/backend-restore.err.log

[Install]
WantedBy=multi-user.target